# private_dot_config/systemd/user/backup-local-check.timer
[Unit]
Description=Weekly local backup integrity check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
