# private_dot_config/systemd/user/backup-offsite-check.service
[Unit]
Description=Offsite backup integrity check
OnFailure=backup-notify@%n.service

[Service]
Type=oneshot
ExecStart=restic check
EnvironmentFile=%h/.config/restic/hetzner-backup.conf
