# private_dot_config/systemd/user/backup-notify@.service
[Unit]
Description=Failure notification for %i

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-notify %i
