# private_dot_config/systemd/user/backup-local.timer
[Unit]
Description=Local backup with restic

[Timer]
OnCalendar=hourly UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
