# private_dot_config/systemd/user/offsite-backup.timer
[Unit]
Description=Offsite backup with restic

[Timer]
OnCalendar=daily UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
